# systemd/user/i3batwarn.service
[Unit]
Description=Run i3 battery warning script

[Service]
ExecStart=/home/roger/Documents/i3batwarn.sh
